# launchd/com.jank.meeting-ai.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jank.meeting-ai.watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jank/code/personal/meeting-ai/scripts/watch-recordings.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Make Homebrew binaries (fswatch, terminal-notifier, ollama) discoverable. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jank/Recordings/Meetings/logs/watch.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jank/Recordings/Meetings/logs/watch.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
