# launchd/com.jank.murmur.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jank.murmur.daemon</string>

  <!-- Entrypoint wrapper: sources config.sh and execs the daemon with stdout/stderr
       redirected to $MEETINGS_BASE/logs, so config.sh is the single source of truth for
       the log location (launchd does no variable expansion in plists). bun is resolved
       from the PATH in EnvironmentVariables below — update that if you upgrade bun. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jank/dev/personal/murmur/launchd/run-daemon.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jank/dev/personal/murmur/src</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- The daemon also prepends these to the PATH it hands children, but set a sane
       PATH so bun + python (asr venv) + ffmpeg + ollama + terminal-notifier all resolve. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/jank/.local/share/mise/installs/bun/latest/bin:/Users/jank/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- No StandardOut/ErrorPath here: run-daemon.sh redirects the daemon's stdout/stderr
       to $MEETINGS_BASE/logs/daemon.{out,err}.log, derived from config.sh. -->


  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
